# slocate.cron
#!/bin/sh

/bin/nice -n 20 /usr/bin/slocate -u -c
